# deploy/systemd/tesla-command-proxy.service.example
[Unit]
Description=Proxy local Tesla pour commandes signées
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=olivier
WorkingDirectory=/home/olivier/git/tesla-charge
ExecStart=/usr/local/bin/tesla-http-proxy -host 127.0.0.1 -port 4443 -cert /home/olivier/git/tesla-charge/raspberry/proxy/tls-cert.pem -tls-key /home/olivier/git/tesla-charge/raspberry/proxy/tls-key.pem -key-file /home/olivier/git/tesla-charge/raspberry/proxy/fleet-key.pem
Restart=always
RestartSec=10
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
